# supervisor/chromes.ini
[program:9225__chrome_regular_pricing]
user=newuser
directory=//project
stdout_logfile=/logs/supervisor/chrome.9225__chrome_regular_pricing.log
stderr_logfile=/logs/supervisor/chrome.9225__chrome_regular_pricing.log
environment=DISPLAY=:90
autorestart=true
command=google-chrome --remote-debugging-port=9225 --user-data-dir=/tmp/pricing --proxy-server="23.249.230.120:65432"

[program:9227__chrome_regular_ranking]
user=newuser
directory=/project
stdout_logfile=/logs/supervisor/chrome.9227__chrome_regular_ranking.log
stderr_logfile=/logs/supervisor/chrome.9227__chrome_regular_ranking.log
environment=DISPLAY=:89
autorestart=true
command=google-chrome --remote-debugging-port=9227 --user-data-dir=/tmp/ranking --proxy-server="23.249.230.120:65432"
