# supervisor/utils.ini
[program:utils__x_display_pricing]
user=newuser
directory=/project
stdout_logfile=/logs/supervisor/chrome.utils__x_display_pricing.log
stderr_logfile=/logs/supervisor/chrome.utils__x_display_pricing.log
autorestart=true
command=Xvfb :90 -screen 0 1536x960x16 -nolisten tcp

[program:utils__x_vnc_pricing]
user=newuser
directory=/project
stdout_logfile=/logs/supervisor/chrome.utils__x_vnc_pricing.log
stderr_logfile=/logs/supervisor/chrome.utils__x_vnc_pricing.log
autorestart=true
command=x11vnc -display :90

[program:utils__x_display_ranking]
user=newuser
directory=/project
stdout_logfile=/logs/supervisor/chrome.utils__x_display_ranking.log
stderr_logfile=/logs/supervisor/chrome.utils__x_display_ranking.log
autorestart=true
command=Xvfb :89 -screen 0 960x720x8 -nolisten tcp

[program:utils__x_vnc_ranking]
user=newuser
directory=/project
stdout_logfile=/logs/supervisor/chrome.utils__x_vnc_ranking.log
stderr_logfile=/logs/supervisor/chrome.utils__x_vnc_ranking.log
autorestart=true
command=x11vnc -display :89
